# app-misc/flatpak-timers/files/system-flatpak-update@.timer
[Unit]
Description=Update Flatpak non-default system-installation

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
