# app-misc/flatpak-timers/files/user-flatpak-update.timer
[Unit]
Description=Update Flatpak user-installation

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
